# reboot-service.service
[Unit]
Description=Inpatient Display Reboot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/inpatient-display
ExecStart=/usr/bin/node /usr/local/bin/inpatient-reboot-service
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=REBOOT_API_KEY_FILE=/etc/inpatient-display/api-key

[Install]
WantedBy=multi-user.target 